# samba/etc/cron.d/krb5_host_ticket_renew
0 5,17 * * * root [ -s /etc/krb5.keytab ] && /usr/bin/kinit -k


